# Docker/centOs/etc/profile.d/java.sh
#! /bin/bash

#
# set environment variables
#
export JAVA_HOME=/opt/jdk1.6.0_45
PATH=$JAVA_HOME/bin:$PATH
