# Docker/centOs/etc/profile.d/start_tomcat.sh
#! /bin/bash

#
# start tomcat
#
service tomcat-1 start
service tomcat-2 start